# config/files/etc/systemd/system/install-1password.service
[Unit]
Description=Install 1Password from Flathub
After=network-online.target
Wants=network-online.target
ConditionPathExists=!/var/lib/flatpak/app/com.1password.1Password

[Service]
Type=oneshot
ExecStart=/usr/bin/flatpak remote-add --system --if-not-exists flathub https://dl.flathub.org/repo/flathub.flatpakrepo
ExecStart=/usr/bin/flatpak install --system --noninteractive flathub com.1password.1Password
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
