# config/files/etc/systemd/user/kanshi.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/kanshi -c /etc/kanshi/config
